# retrograde/release/retrograde.service
[Unit]
# Human readable name of the unit
Description=Retrograde Display service

[Install]
WantedBy=default.target

[Service]
Restart=on-failure
